# root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
export APP_KEY=$(cat /config/SNIPE_IT_APP_KEY.txt)
exec /usr/sbin/php-fpm7 -F
