# root/etc/crontabs/abc
# min   hour    day     month   weekday command
*       *       *       *       *       php /app/www/artisan schedule:run >> /dev/null 2>&1
