# root/etc/s6-overlay/s6-rc.d/init-snipe-it-config/up
/etc/s6-overlay/s6-rc.d/init-snipe-it-config/run
