# root/etc/s6-overlay/s6-rc.d/svc-php-fpm/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

export APP_KEY=$(cat /config/SNIPE_IT_APP_KEY.txt)

exec /usr/sbin/php-fpm82 -F
